# app-admin/istatd/files/istatd-init
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

PID="/run/istatd/istat.pid"

checkTmp() {
    cache_dir="/var/cache/istat"
    run_dir="/run/istatd"

    if [ ! -d "${cache_dir}" ]; then
        mkdir -p "${cache_dir}"
        chown istatd:istatd "${cache_dir}"
    fi

    if [ ! -d "${run_dir}" ]; then
        mkdir -p "${run_dir}"
        chown istatd:istatd "${run_dir}"
    fi
}

depend() {
    need net
}

start() {
	ebegin "Starting istatd"
    checkTmp
	start-stop-daemon --start \
		--exec /usr/bin/istatd \
		--make-pidfile \
		--chuid istatd:istatd \
		--pidfile ${PID} \
		--background \
		-- "${ISTATD_OPTS}" \
		>> /var/log/istatd 2>&1
	eend $?
}

stop() {
	ebegin "Stopping istatd"
	start-stop-daemon --stop \
		--pidfile ${PID}
	eend $?
}

